# scripts/com.james.pool-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.james.pool-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>REPLACE_WITH_VENV_PYTHON</string>
        <!-- e.g. /Users/james/pool-monitor/.venv/bin/python3 -->
        <string>-m</string>
        <string>monitor.runner</string>
    </array>

    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_PROJECT_PATH</string>
    <!-- e.g. /Users/james/pool-monitor -->

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>REPLACE_WITH_PROJECT_PATH/logs/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_PROJECT_PATH/logs/launchd.log</string>
</dict>
</plist>
